# infrastructure/golden-image/services/novnc.service
[Unit]
Description=noVNC Web Interface for Golden Image Desktop
After=vnc-server.service
Requires=vnc-server.service

[Service]
Type=simple
ExecStart=/usr/bin/websockify --web=/usr/share/novnc 6080 localhost:5901
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
